# crontab
# Proactive messaging cron jobs
# Format: minute hour day month weekday command
#
# T-1 day reminder: run daily at 16:00 UTC
# Day-of reminder: run daily at 16:00 UTC (same job handles both)
0 16 * * * curl -fsS -X POST http://127.0.0.1:3000/internal/proactive/check-reminders || true

# Weekly announcement: run every Monday at 15:00 UTC
0 18 * * 1 curl -fsS -X POST http://127.0.0.1:3000/internal/proactive/weekly-announcement || true
